# templates/druid-overlord.service.erb
[Unit]
Description=Druid Overlord Node

[Service]
Type=simple
WorkingDirectory=/etc/druid/overlord/
ExecStart=/usr/bin/java <%= @jvm_opts.join(' ') %> -classpath .:<%= @install_dir %>/druid/lib/* io.druid.cli.Main server overlord
SuccessExitStatus=130 143
Restart=on-failure

[Install]
WantedBy=multi-user.target
