# templates/druid-historical.service.erb
[Unit]
Description=Druid Historical Node

[Service]
Type=simple
WorkingDirectory=/etc/druid/historical/
ExecStart=/usr/bin/java <%= @jvm_opts.join(' ') %> -classpath .:/usr/local/lib/druid/lib/* io.druid.cli.Main server historical
SuccessExitStatus=130 143
Restart=on-failure

[Install]
WantedBy=multi-user.target
