# templates/druid-broker.service.erb
[Unit]
Description=Druid Broker Node

[Service]
Type=simple
WorkingDirectory=/etc/druid/broker/
ExecStart=/usr/bin/java <%= @jvm_opts.join(' ') %> -classpath .:<%= @install_dir %>/druid/lib/* io.druid.cli.Main server broker
SuccessExitStatus=130 143
Restart=on-failure

[Install]
WantedBy=multi-user.target
