# templates/druid-realtime.service.erb
[Unit]
Description=Druid Realtime Node

[Service]
Type=simple
WorkingDirectory=/etc/druid/realtime/
ExecStart=/usr/bin/java <%= @jvm_opts.join(' ') %> -classpath .:<%= scope.lookupvar("::#{@module_name}::install_dir") %>/druid/lib/* io.druid.cli.Main server realtime
SuccessExitStatus=130 143
Restart=on-failure

[Install]
WantedBy=multi-user.target
