# templates/druid-middle_manager.service.erb
[Unit]
Description=Druid Middle Manager Node

[Service]
Type=simple
WorkingDirectory=/etc/druid/middle_manager/
ExecStart=/usr/bin/java <%= @jvm_opts.join(' ') %> -classpath .:<%= scope.lookupvar("::#{@module_name}::install_dir") %>/druid/lib/* io.druid.cli.Main server middleManager
SuccessExitStatus=130 143
Restart=on-failure

[Install]
WantedBy=multi-user.target
